# lekcja-08-bash/backup/backup.service
[Unit]
Description=Auto backup

[Service]
ExecStart=/home/vboxuser/backup/backup.sh

[Install]
WantedBy=default.target
